# deploy/systemd/var-lit-v1-runtime.service.example
[Unit]
Description=Var-Lit V1 trading runtime
Requires=var-lit-v1-chrome.service
After=network-online.target var-lit-v1-chrome.service
Wants=network-online.target

[Service]
Type=simple
User=varlit
Group=varlit
WorkingDirectory=/opt/var-lit-v1
Environment=VARIATIONAL_ENV_FILE=/var/lib/var-lit-v1/runtime.env
ExecStartPre=/usr/bin/test -f /var/lib/var-lit-v1/runtime.env
ExecStart=/opt/var-lit-v1/deploy/run_runtime.sh
Restart=on-failure
RestartSec=3
TimeoutStopSec=40
KillSignal=SIGTERM
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
