# deploy/systemd/var-lit-v1-display.service.example
[Unit]
Description=Var-Lit V1 virtual display
After=network.target

[Service]
Type=simple
User=varlit
Group=varlit
ExecStart=/usr/bin/Xvfb :99 -screen 0 1280x900x24 -nolisten tcp -noreset
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
